# tmp.web.server/consul.service
[Unit]
Description=Consul client service

[Service]
Type=simple
ExecStart=/bin/bash /tmp/web.server/start_consul.sh

[Install]
