# tmp.web.server/web.service
[Unit]
Description=Web Service

[Service]
Type=simple
ExecStart=/bin/bash /tmp/web.server/start_web.sh

[Install]
