# system/consul_server.service
[Unit]
Description=Consul server service

[Service]
Type=simple
ExecStart=/bin/bash /vagrant/scripts/start_consul.sh

[Install]
