# tmp.consul.server/consul.service
[Unit]
Description=Consul server service

[Service]
Type=simple
ExecStart=/bin/bash /tmp/consul.server/scripts/start_consul.sh

[Install]
